# autostart/autostart.service
[Unit]
Description=XDG Autostart (via dex)
Wants=basic.target

[Service]
Type=oneshot
ExecStart=dex -a

[Install]
WantedBy=default.target
