# aconfmgr/files/etc/systemd/system/docker.service.d/resourcecontrol.conf
[Service]
AllowedCPUs=1-3
MemoryHigh=4G
MemoryMax=6G
