# rofi/rofi-drun@.service
[Service]
Type=forking
ExecStart=/usr/bin/rofi -show drun
