# tlp/netlink-change.path
[Path]
PathChanged=/run/systemd/netif/state
Unit=netlink-change.service

[Install]
WantedBy=default.target
