# polybar/polybar@one.service
[Unit]
Description=Polybar %I

[Service]
ExecStart=/usr/bin/polybar -r %i
Restart=always

[Install]
WantedBy=default.target
