# pass/pass-sync.service
[Unit]
Description=`pass` db sync
After=default.target
# Only try once an hour
# (otherwise repeated gpg-agent pinentry dialogues)
StartLimitBurst=1
StartLimitIntervalSec=3600

[Service]
Type=oneshot
ExecStart=/usr/bin/env pass-sync --down
Restart=on-failure
RemainAfterExit=yes
ExecStop=/usr/bin/env pass-sync --up

[Install]
WantedBy=default.target
