# Keybase/keybase.target
[Unit]
Description=Keybase
Wants=keybase.service kbfs.service

[Install]
WantedBy=default.target
